# systemd/borg-test-restore.service
[Unit]
Description=Monthly BorgBackup restore verification
Documentation=https://github.com/pavlojs/flex-backup-system
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/root/.backup-secrets.env
ExecStart=/root/borg-test-restore.sh

# Resource limits
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUQuota=50%

# Security hardening
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=false
ReadWritePaths=/var/log /var/backups /tmp /root

[Install]
WantedBy=multi-user.target
