# systemd/restic-backup.service
[Unit]
Description=Restic Backup to Cloudflare R2
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/root/.backup-secrets.env
ExecStart=/root/restic-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUQuota=50%

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log /var/lib/docker/volumes /home/pavlojs/apps /root
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
